# .openclaw/workspace/projects/Axe-Capital/ops/systemd/axe-monthly-specialists.timer
[Unit]
Description=Run Axe Capital portfolio specialists monthly

[Timer]
OnCalendar=monthly
Persistent=true
Unit=axe-monthly-specialists.service

[Install]
WantedBy=timers.target
